# docker/supervisord-gpu-ssl.conf
[supervisord]
nodaemon=true
user=omega
logfile=/app/logs/supervisord-gpu.log
pidfile=/app/tmp/supervisord-gpu.pid
childlogdir=/app/logs

[program:gpu-ssl-handler]
command=/app/gpu-ssl-handler.sh
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/gpu-ssl-handler.log
stderr_logfile=/app/logs/gpu-ssl-handler.log
priority=100

[program:gpu-service]
command=python -m gunicorn gpu_service:app --bind 0.0.0.0:8001 --workers 2 --worker-class sync --timeout 600 --keepalive 30 --max-requests 1000 --preload
directory=/app
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/gpu-service.log
stderr_logfile=/app/logs/gpu-service.log
environment=
    PYTHONPATH="/app",
    CUDA_VISIBLE_DEVICES="0",
    NVIDIA_VISIBLE_DEVICES="all",
    SSL_ENABLED="true",
    GPU_MEMORY_STRATEGY="balanced"
priority=200
depends_on=gpu-ssl-handler

[program:gpu-monitor]
command=python /app/scripts/gpu_monitor.py
directory=/app
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/gpu-monitor.log
stderr_logfile=/app/logs/gpu-monitor.log
priority=300

[program:model-loader]
command=python /app/scripts/model_preloader.py
directory=/app
autostart=true
autorestart=false
startsecs=10
user=omega
stdout_logfile=/app/logs/model-loader.log
stderr_logfile=/app/logs/model-loader.log
priority=150

[group:omega-gpu]
programs=gpu-ssl-handler,gpu-service,gpu-monitor,model-loader