# docker/supervisord-ssl.conf
[supervisord]
nodaemon=true
user=omega
logfile=/app/logs/supervisord.log
pidfile=/app/tmp/supervisord.pid
childlogdir=/app/logs

[program:ssl-handler]
command=/app/ssl-handler.sh
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/ssl-handler.log
stderr_logfile=/app/logs/ssl-handler.log
priority=100

[program:gunicorn]
command=gunicorn main:app --bind 0.0.0.0:8000 --workers 4 --worker-class gevent --worker-connections 1000 --timeout 300 --keepalive 30
directory=/app
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/gunicorn.log
stderr_logfile=/app/logs/gunicorn.log
environment=
    PYTHONPATH="/app",
    SSL_ENABLED="true",
    SECURITY_HEADERS="true"
priority=200

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
user=root
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.log
priority=300
depends_on=ssl-handler,gunicorn

[program:ssl-monitor]
command=python /app/scripts/ssl_cert_manager.py monitor
directory=/app
autostart=true
autorestart=true
user=omega
stdout_logfile=/app/logs/ssl-monitor.log
stderr_logfile=/app/logs/ssl-monitor.log
priority=400

[group:omega-api]
programs=ssl-handler,gunicorn,nginx,ssl-monitor